# smartpanel.service
[Unit]
Description=Smart Panel Service
After=network.target

[Service]
WorkingDirectory=/home/pi/Apps/smart-panel
ExecStart=/bin/bash run.sh
Restart=on-failure
User=pi
Environment=PORT=8000

[Install]
WantedBy=multi-user.target